# systemd/yt-dashboard.service
[Unit]
# TEMPLATE — placeholders sao substituidos pelo scripts/setup-canal:
#   __INSTANCE_DIR__   ex: /home/user/projetos/yt-pub-lives7
#   __INSTANCE_NAME__  ex: yt-pub-lives7
#   __PORT__           ex: 8097
Description=YT Pub Lives __INSTANCE_NAME__ Dashboard
After=network.target

[Service]
Type=simple
WorkingDirectory=__INSTANCE_DIR__/dashboard
ExecStart=/usr/bin/python3 __INSTANCE_DIR__/dashboard/server.py __PORT__
Restart=always
RestartSec=5
Environment=PATH=/usr/bin:%h/.local/bin:%h/.npm-global/bin
Environment=GWS_CONFIG_DIR=__INSTANCE_DIR__/config

[Install]
WantedBy=default.target
