# master-dashboard/yt-master-dashboard.service
[Unit]
Description=YT Pub Master Dashboard (all instances)
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/nmaldaner/projetos/yt-pub-lives2/master-dashboard
ExecStart=/usr/bin/python3 /home/nmaldaner/projetos/yt-pub-lives2/master-dashboard/server.py
Restart=always
RestartSec=10
Environment=PATH=/usr/bin:/home/nmaldaner/.local/bin
Environment=GWS_CONFIG_DIR=/home/nmaldaner/projetos/yt-pub-lives2/config
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
